# ops/postdoctor-frontend.service
[Unit]
Description=PostDoctor Frontend (Static Server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/postdoctor/client
Environment=FRONTEND_ROOT=/opt/postdoctor/client/current
Environment=FRONTEND_HOST=127.0.0.1
Environment=FRONTEND_PORT=4173
ExecStart=/usr/bin/node /opt/postdoctor/client/scripts/static-server.cjs
Restart=always
RestartSec=5
User=postdoctor
Group=postdoctor

[Install]
WantedBy=multi-user.target
